# systemd/miner-exclusion.service
[Unit]
Description=MINER — Moniteur liste d'exclusion (nettoyage temps réel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Projets_en_Production/MINER
EnvironmentFile=/home/ubuntu/Projets_en_Production/MINER/.env
ExecStart=/home/ubuntu/Projets_en_Production/MINER/.venv/bin/python monitor_exclusion.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/var/log/miner/monitor_exclusion.log
StandardError=append:/var/log/miner/monitor_exclusion.log

[Install]
WantedBy=multi-user.target
